# deploy/qk80-keepalive.timer
[Unit]
Description=qk 80B keepalive ping every 3 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=3min
AccuracySec=30s

[Install]
WantedBy=timers.target
