# deploy/qk80-keepalive.service
[Unit]
# Keep the standalone 80B's Metal residency warm.
#
# The weights are mlocked on midnight, but the FIRST request after an idle gap
# still runs at 5-20 tok/s while Metal re-warms its residency — an agent's first
# turn after a break waits ~11 s for what should take 1 s. A 1-token request
# every 3 minutes keeps it warm; it costs one decode step and never touches a
# slot for long enough to matter (2 slots, and the driver frees it immediately).
Description=qk 80B keepalive ping (keeps Metal residency warm)

[Service]
Type=oneshot
ExecStart=/usr/bin/curl -sS -m 30 -o /dev/null \
  -H "content-type: application/json" \
  -d "{\"model\":\"qk\",\"max_tokens\":1,\"temperature\":0,\"messages\":[{\"role\":\"user\",\"content\":\"ping\"}]}" \
  http://REDACTED-CLUSTER-IP:8080/v1/messages
